# .claude/crontab
# Claude Audit Cron Jobs
# Split Lease - Scheduled Audit Commands
# Format: minute hour day month dayofweek command

# Run audit-twilio-sms-mocking at 10:15 AM daily
15 10 * * * /home/splitlease/Documents/splitlease/.claude/cron-audits.sh audit-twilio-sms-mocking

# Run audit-vitest-rtl-setup at 11:15 AM daily
15 11 * * * /home/splitlease/Documents/splitlease/.claude/cron-audits.sh audit-vitest-rtl-setup

# Run audit-custom-hook-tests at 11:45 AM daily
45 11 * * * /home/splitlease/Documents/splitlease/.claude/cron-audits.sh audit-custom-hook-tests

# Run audit-barrel-files at 12:15 PM daily
15 12 * * * /home/splitlease/Documents/splitlease/.claude/cron-audits.sh audit-barrel-files
